# my-claw/systemd/cypherclaw-observer-ollama.service
[Unit]
Description=CypherClaw Observer Ollama
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=simple
User=ollama
Group=ollama
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="OLLAMA_HOST=127.0.0.1:11435"
Environment="OLLAMA_MODELS=/usr/share/ollama/.ollama/models"
Environment="TMPDIR=/run/cypherclaw-tmp"
Environment="OLLAMA_MAX_LOADED_MODELS=1"
Environment="OLLAMA_KEEP_ALIVE=15m"
Environment="OLLAMA_FLASH_ATTENTION=1"
Environment="OLLAMA_KV_CACHE_TYPE=q8_0"
Environment="OLLAMA_NUM_PARALLEL=1"
ExecStartPre=+/usr/bin/install -d -m 1777 /run/cypherclaw-tmp
ExecStart=/usr/local/bin/ollama serve
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
