# my-claw/systemd/cypherclaw-sdp-runner.service
[Unit]
Description=CypherClaw SDP Runner
After=cypherclaw-bootstrap.service cypherclaw.service network-online.target
Requires=cypherclaw-bootstrap.service
ConditionPathExists=/home/user/cypherclaw

[Service]
Type=simple
User=user
WorkingDirectory=/run/cypherclaw-tmp/workdir/cypherclaw-work
ExecStart=/home/user/cypherclaw/tools/sdp_runner_launcher.sh
Restart=on-failure
RestartSec=10
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
