# my-claw/systemd/cypherclaw-idyom-kl-audit.timer
[Unit]
Description=CypherClaw Weekly IDyOM LTM KL-Divergence Audit Timer (T-031)
Requires=cypherclaw-idyom-kl-audit.service

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=15min
Unit=cypherclaw-idyom-kl-audit.service

[Install]
WantedBy=timers.target
