# my-claw/systemd/cypherclaw-bootstrap.service
[Unit]
Description=CypherClaw Bootstrap
After=local-fs.target network.target
Before=cypherclaw.service cypherclaw-sdp-runner.service
ConditionPathExists=/home/user/cypherclaw

[Service]
Type=oneshot
User=user
WorkingDirectory=/home/user/cypherclaw
ExecStart=/home/user/cypherclaw/tools/init_workdir.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
