# my-claw/tools/boot_hardening.service
[Unit]
Description=CypherClaw Boot Hardening
After=network.target redis-server.service postgresql.service
Before=cypherclaw.service

[Service]
Type=oneshot
ExecStart=/home/user/cypherclaw/tools/boot_hardening.sh
RemainAfterExit=yes
User=root

[Install]
WantedBy=multi-user.target
