# my-claw/tools/cypherclaw-dashboard.timer
[Unit]
Description=Run CypherClaw Dashboard Generator every 60 seconds

[Timer]
OnBootSec=1min
OnUnitActiveSec=60s
Unit=cypherclaw-dashboard.service

[Install]
WantedBy=timers.target
